# containers/mpd/mpd.container
# Podman Quadlet file for MPD
# Place in ~/.config/containers/systemd/mpd.container
# Then run: systemctl --user daemon-reload && systemctl --user start mpd

[Unit]
Description=Music Player Daemon (containerized)
After=pipewire.service pipewire-pulse.service
Wants=pipewire-pulse.service

[Container]
Image=localhost/mpd:latest
ContainerName=mpd
UserNS=keep-id

# Music library - mount your music directory
Volume=%h/music:/var/lib/mpd/music:ro

# Persistent data (playlists, database, state)
Volume=%h/.local/share/mpd:/var/lib/mpd/data
Volume=%h/.config/mpd/playlists:/var/lib/mpd/playlists

# PipeWire/PulseAudio socket for audio output
Volume=/run/user/%U/pulse:/run/user/1000/pulse
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native

# Network
PublishPort=6600:6600

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
