# speech-engines/systemd/chatterbox-tts.service
[Unit]
Description=Chatterbox TTS Server (OpenAI-compatible, ROCm GPU)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/src/1st-level/@rag/speech-engines/chatterbox-server
EnvironmentFile=%h/src/1st-level/@rag/speech-engines/config/chatterbox.env
Environment=PATH=/opt/rocm/bin:/usr/bin:/usr/sbin
ExecStart=%h/src/1st-level/@rag/speech-engines/.venv-chatterbox/bin/python server.py
Restart=on-failure
RestartSec=10

# Resource limits
MemoryMax=16G
TimeoutStartSec=120

[Install]
WantedBy=default.target
