# speech-engines/systemd/whisper-stt.service
[Unit]
Description=Whisper STT Server (OpenAI-compatible, whisper.cpp)
After=network.target

[Service]
Type=simple
ExecStart=%h/src/salt/speech-engines/.venv-piper/bin/python %h/src/salt/speech-engines/whisper-server.py \
    --port 8002 \
    --host 0.0.0.0 \
    --model %h/src/salt/speech-engines/voices/ggml-large-v3-turbo.bin \
    --whisper-bin whisper-cli \
    --ld-library-path %h/src/salt/speech-engines/whisper.cpp/build/src:%h/src/salt/speech-engines/whisper.cpp/build/ggml/src:%h/src/salt/speech-engines/whisper.cpp/build/ggml/src/ggml-hip
Restart=on-failure
RestartSec=5

MemoryMax=4G

[Install]
WantedBy=default.target
