# states/units/llama_embed.container
{#- Quadlet container unit template for llama.cpp embedding server
    (feature 087-containerize-services). Rendered by container_service macro.
    Context:
      catalog_entry:    dict from service_catalog.yaml[llama_embed]
      image:            fully-resolved registry/image@sha256:... reference
      expanded_mounts:  list of {host, container, mode}
      user_scope:       False (system scope per FR-017)

    The upstream ghcr.io/ggml-org/llama.cpp:server-vulkan image has
    llama-server as its entrypoint and accepts CLI args passed via Quadlet's
    Exec= key. Model path is /models/<file> because the bind-mount maps
    the host model directory to /models inside the container.
-#}
{%- import_yaml 'data/llama_embed.yaml' as embed -%}
[Unit]
Description=llama.cpp embedding server (containerized, Qwen3-Embedding-8B, Vulkan)
Documentation=https://github.com/ggml-org/llama.cpp
After=network-online.target local-fs.target mnt-one.mount
Wants=network-online.target
Requires=mnt-one.mount

[Container]
Image={{ image }}
ContainerName=llama_embed
PublishPort=127.0.0.1:{{ catalog_entry.port }}:{{ catalog_entry.port }}
{% for m in expanded_mounts -%}
Volume={{ m.host }}:{{ m.container }}:{{ m.mode }}
{% endfor -%}
AddDevice=/dev/dri/renderD128
# Numeric GID at render time — see ollama.container for the rationale.
GroupAdd={{ salt['cmd.run_stdout']("getent group render | cut -d: -f3", python_shell=True) }}
# Override the image entrypoint args — the base image runs llama-server
# but needs explicit configuration for the embedding workload.
Exec=--model /models/{{ embed.file }} --embeddings --pooling {{ embed.pooling }} --ctx-size {{ embed.context }} --gpu-layers {{ embed.gpu_layers }} --host 0.0.0.0 --port {{ catalog_entry.port }}
HealthCmd=curl -sf http://127.0.0.1:{{ catalog_entry.port }}{{ catalog_entry.health_path }} || exit 1
# 2s poll interval + 60s start grace:
# llama-server needs ~2s to initialize Vulkan + mmap the 6 GB model + bind
# the port, so the first probe at container-create reliably fails. With
# the podman default HealthInterval=30s, Notify=healthy only fires on the
# second probe at t+30s, turning every restart into a 30.47s cold-start
# number (pure probe-cadence artifact, not a workload regression). A 2s
# interval lets the second probe hit at ~t+2s — close to real readiness.
# HealthStartPeriod=60s gives the container a grace window where failing
# probes do not count toward the retry budget, tolerating slow Vulkan
# device init on a cold GPU without flipping the unit to failed.
HealthInterval=2s
HealthStartPeriod=60s
HealthRetries=3
HealthTimeout=5s
Notify=healthy

[Service]
# Fail-loud GPU preflight check (FR-007, research Decision 2). Vulkan-only:
# no /dev/kfd required. llama-server also refuses to start without a Vulkan
# device at the binary level, so fail-loud is double-guaranteed.
ExecStartPre=/usr/bin/test -e /dev/dri/renderD128
# Bounded restart policy per FR-015.
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
# TimeoutStartSec=180: Vulkan device init + Q5_K_M 8B model load on cold GPU
# can take up to 90s; leave headroom.
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
