# states/units/ollama.service
[Unit]
Description=Ollama LLM Server
After=network-online.target

[Service]
ExecStart=/usr/bin/ollama serve
User={{ user }}
Group={{ user }}
Restart=always
RestartSec=3
WorkingDirectory={{ home }}
Environment="HOME={{ home }}"
Environment="OLLAMA_HOST=127.0.0.1:{{ ollama_port }}"
Environment="OLLAMA_MODELS={{ mnt_one }}/ollama/models"
Environment="OLLAMA_VULKAN=1"
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=yes
# PrivateDevices=yes — disabled: Vulkan backend needs /dev/dri for GPU inference
ReadWritePaths={{ mnt_one }}/ollama
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
