# states/units/salt-daemon.service.j2
[Unit]
Description=Salt state daemon (pre-loaded SMinion for fast state.apply)
Documentation=file://{{ project_dir }}/scripts/salt-daemon.py
After=network.target

[Service]
Type=simple
ExecStart={{ project_dir }}/.venv/bin/python3 {{ project_dir }}/scripts/salt-daemon.py \
    --config-dir {{ project_dir }}/.salt_runtime \
    --log-level warning
WorkingDirectory={{ project_dir }}
Restart=on-failure
RestartSec=5
Environment=XDG_RUNTIME_DIR={{ runtime_dir }}
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path={{ runtime_dir }}/bus
# State application requires root
User=root
# Clean up socket on stop
ExecStopPost=/bin/rm -f /run/salt-daemon.sock
PrivateTmp=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
