# states/units/salt-daemon.service
[Unit]
Description=Salt state daemon (pre-loaded SMinion for fast state.apply)
Documentation=file:///home/neg/src/salt/scripts/salt-daemon.py
After=network.target

[Service]
Type=simple
ExecStart=/home/neg/src/salt/.venv/bin/python3 /home/neg/src/salt/scripts/salt-daemon.py \
    --config-dir /home/neg/src/salt/.salt_runtime \
    --socket /tmp/salt-daemon.sock \
    --log-level warning
WorkingDirectory=/home/neg/src/salt
Restart=on-failure
RestartSec=5
# State application requires root
User=root
# Clean up socket on stop
ExecStopPost=/bin/rm -f /tmp/salt-daemon.sock

[Install]
WantedBy=multi-user.target
