# states/units/user/gpu-unstick.service
[Unit]
Description=Reset stuck amdgpu PM counters (RDNA3 workaround)

[Service]
Type=oneshot
ExecStart=%h/.local/bin/gpu-unstick
