# states/units/user/gpu-unstick.timer
[Unit]
Description=Periodically reset stuck amdgpu PM counters

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
