# states/units/user/nanoclaw.service
[Unit]
Description=NanoClaw AI Agent
After=network-online.target proxypilot.service
Wants=proxypilot.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
WorkingDirectory=%h/.local/share/nanoclaw
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=10

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=tmpfs
BindPaths=%h/.local/share/nanoclaw:%h/.local/share/nanoclaw
BindReadOnlyPaths=%h/.config/nanoclaw:%h/.config/nanoclaw
BindPaths=%h/.local/share/nanoclaw/store:%h/.local/share/nanoclaw/store
BindPaths=%h/.local/share/nanoclaw/data:%h/.local/share/nanoclaw/data
BindPaths=%h/.local/share/nanoclaw/groups:%h/.local/share/nanoclaw/groups
CapabilityBoundingSet=
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
ProtectClock=true
ProtectHostname=true
SystemCallArchitectures=native
# V8 JIT requires W+X pages — cannot use MemoryDenyWriteExecute
# MemoryDenyWriteExecute=true

# Environment
Environment=NODE_ENV=production

[Install]
WantedBy=default.target
