# states/units/user/openclaw-gateway.service
[Unit]
Description=OpenClaw AI agent gateway
Documentation=https://docs.openclaw.ai
After=network-online.target proxypilot.service
Wants=network-online.target
Requires=proxypilot.service
ConditionPathExists=%h/.openclaw/openclaw.json

[Service]
ExecStart=%h/.local/bin/openclaw gateway
WorkingDirectory=%h/.openclaw
SyslogIdentifier=openclaw-gateway
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateDevices=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=default.target
