# states/units/user/openclaw-health.service
[Unit]
Description=OpenClaw gateway health check
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/openclaw-health-check
TimeoutStartSec=60
