# states/units/user/opencode-telegram-bot.service
[Unit]
Description=OpenCode Telegram Bot — mobile client for OpenCode
After=network-online.target opencode-serve.service
Wants=network-online.target
Requires=opencode-serve.service
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
ExecStart=%h/.local/bin/opencode-telegram start
WorkingDirectory=%h
SyslogIdentifier=opencode-telegram-bot
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
MemoryMax=512M
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
EnvironmentFile=%h/.config/opencode-telegram-bot/.env
NoNewPrivileges=true
CapabilityBoundingSet=
SystemCallArchitectures=native
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
PrivateTmp=yes
KeyringMode=private
UMask=0077

[Install]
WantedBy=default.target
