# states/units/user/opencode-telegram-bot.container
{#- Quadlet container unit template for OpenCode Telegram bot
    (feature 087-containerize-services, US3 structural scaffolding).
    DIGEST DEFERRED per research Decision 7 — upstream image pending.
    Context: see telethon-bridge.container.
-#}
[Unit]
Description=OpenCode Telegram bot (containerized)
After=default.target opencode-serve.service
Requires=opencode-serve.service

[Container]
Image={{ image }}
ContainerName=opencode-telegram-bot
{% for m in expanded_mounts -%}
Volume={{ m.host }}:{{ m.container }}:{{ m.mode }}
{% endfor -%}
{% if catalog_entry.env_file -%}
EnvironmentFile={{ catalog_entry.env_file }}
{% endif -%}
HealthCmd=pgrep -f opencode-telegram || exit 1
HealthInterval=30s
HealthRetries=3
HealthTimeout=10s
Notify=healthy

[Service]
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
TimeoutStartSec=60

[Install]
WantedBy=default.target
