# states/units/user/openclaw-gateway.service
[Unit]
Description=OpenClaw AI agent gateway
Documentation=https://docs.openclaw.ai
After=network-online.target proxypilot.service
Wants=network-online.target proxypilot.service
ConditionPathExists=%h/.openclaw/openclaw.json
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
ExecStartPre=%h/.local/bin/openclaw-sanitize-config
ExecStart=%h/.local/bin/openclaw gateway
WorkingDirectory=%h/.openclaw
SyslogIdentifier=openclaw-gateway
LogFilterPatterns=~bot[0-9]+:AA
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
MemoryMax=2G
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
NoNewPrivileges=true
# PrivateDevices=yes — disabled: mpv needs /dev/dri for GPU video output
# MemoryDenyWriteExecute=yes — disabled: V8 JIT requires W+X pages
CapabilityBoundingSet=
SystemCallArchitectures=native
SystemCallFilter=@system-service
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ProtectSystem=strict
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectHostname=yes
PrivateTmp=yes
KeyringMode=private
UMask=0077

[Install]
WantedBy=default.target
