# states/units/user/salt-drift.service
[Unit]
Description=Salt package drift detection

[Service]
Type=oneshot
ExecStart=%h/src/salt/scripts/drift-notify.sh
WorkingDirectory=%h/src/salt
