# states/units/user/salt-monitor-watchdog.timer
[Unit]
Description=Salt monitor watchdog timer

[Timer]
OnBootSec=120
OnUnitActiveSec=60
Persistent=true

[Install]
WantedBy=timers.target
