# states/units/user/vdirsyncer.timer
[Unit]
Description=Synchronize calendars and contacts timer
[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target
